# deploy/latasya-backup.timer
[Unit]
Description=Weekly backup of Latasya ERP database

[Timer]
# Sunday 19:00 UTC = Monday 02:00 Jakarta — overnight after the business week.
OnCalendar=Sun *-*-* 19:00:00 UTC
# Run a missed firing as soon as the host is back (covers reboots / downtime
# spanning the schedule).
Persistent=true
Unit=latasya-backup.service

[Install]
WantedBy=timers.target
